# etc/supervisord.d/novnc.ini
[program:websockify]
priority = 40
command = /usr/bin/websockify --ssl-only --web=/usr/share/novnc/ --cert=/etc/pki/tls/certs/novnc.pem 8080 localhost:5901
stdout_logfile=/var/log/supervisor/novnc.log
redirect_stderr=true