# etc/supervisord.d/vncserver.ini
[program:Xvnc]
priority = 10
command = /usr/bin/vncserver -fg
user = USERNAME
environment = HOME="/home/USERNAME",USER="USERNAME",LANG="en_US.UTF-8"
stdout_logfile=/var/log/supervisor/vncserver.log
redirect_stderr=true
autorestart=true
