# etc/supervisord.d/xrdp.ini
[program:xrdp-sesman]
priority = 20
command = /usr/sbin/xrdp-sesman --nodaemon
stdout_logfile=/var/log/supervisor/xrdp-sesman.log
redirect_stderr=true

[program:xrdp]
priority = 30
command = /usr/sbin/xrdp -nodaemon
stdout_logfile=/var/log/supervisor/xrdp.log
redirect_stderr=true